# build/rabbitmq/mquseradd.service
[Unit]
Description=mquseradd
After=network.target rabbitmq-server.service
Requires=rabbitmq-server.service

[Service]
Type=simple
ExecStart=/opt/ARL-ARM64-CentOS8/build/rabbitmq/mqadduser.sh
Restart=no

[Install]
WantedBy=multi-user.target
